# packaging/linux/codex-app-updater.service
[Unit]
Description=Codex App Updater
After=graphical-session.target network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/codex-app-updater daemon
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
