# contrib/user-local-install/files/.config/systemd/user/codex-app-update.service
[Unit]
Description=Update Codex app user-local install

[Service]
Type=oneshot
ExecStart=%h/.local/bin/codex-app-update --quiet
